# deploy/systemd/membridge-agent.service
[Unit]
Description=Membridge Agent Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/membridge
EnvironmentFile=%h/membridge/.env.agent
ExecStart=%h/membridge/.venv/bin/python -m uvicorn agent.main:app --host 0.0.0.0 --port 8001
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=membridge-agent

[Install]
WantedBy=multi-user.target
